# roles/nfs/templates/nfs.automount.j2
[Unit]
Description=Automount NFS for {{ nfs.mountpoint }}
Wants=network-online.target
After=network-online.target

[Automount]
Where={{ nfs.mountpoint }}
TimeoutIdleSec={{ nfs.idle_timeout_sec | default(nfs_idle_timeout_sec) }}

[Install]
WantedBy=multi-user.target
